# deploy/systemd/hiddenoasis-inventory-frontend.service
[Unit]
Description=Hidden Oasis Inventory Frontend
After=network.target hiddenoasis-inventory-backend.service

[Service]
User=hiddenoasis
Group=hiddenoasis
WorkingDirectory=/opt/inventory-procurement-online/frontend
Environment=NODE_ENV=production
Environment=HOSTNAME=127.0.0.1
Environment=PORT=3300
EnvironmentFile=/etc/hiddenoasis/inventory-frontend.env
ExecStart=/usr/bin/node /opt/inventory-procurement-online/frontend/.next/standalone/server.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
