# deploy/systemd/hiddenoasis-inventory-backup.service
[Unit]
Description=Hidden Oasis Inventory Database Backup
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=oneshot
User=hiddenoasis
Group=hiddenoasis
WorkingDirectory=/opt/inventory-procurement-online/backend
EnvironmentFile=/etc/hiddenoasis/inventory-backend.env
ExecStart=/opt/inventory-procurement-online/backend/.venv/bin/python scripts/backup_database.py
