# deploy/systemd/hiddenoasis-inventory-worker.service
[Unit]
Description=Hidden Oasis Inventory Integration Worker
After=network.target postgresql.service hiddenoasis-inventory-backend.service
Wants=postgresql.service

[Service]
User=hiddenoasis
Group=hiddenoasis
WorkingDirectory=/opt/inventory-procurement-online/backend
EnvironmentFile=/etc/hiddenoasis/inventory-backend.env
ExecStart=/bin/sh -c 'while true; do /opt/inventory-procurement-online/backend/.venv/bin/python scripts/process_integration_events.py --limit 25; sleep 10; done'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
